# pipeline/sp-pipeline.cron
# Daily Spurti pipeline at 05:45 UTC = 11:15 IST, right after the 09:00-11:00 IST
# mandatory session. Runs sp-pipeline.sh: #zoomupdate -> sp-rubric-build (APPLY)
# -> sync-spurti-from-sakshi, so interns see the morning session's SP by
# ~11:30 IST instead of next-day. flock keeps it single-instance.
# Install: sudo install -o root -g root -m 644 \
#            /var/samagama/server/sp-pipeline.cron /etc/cron.d/sp-pipeline
# Added 2026-06-04.
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
45 5 * * * samagama /usr/bin/flock -n /tmp/sp-pipeline.lock /var/samagama/server/sp-pipeline.sh >> /var/samagama/server/logs/sp-pipeline.log 2>&1
# Catchup AI summary fetch at 13:00 IST (07:30 UTC) daily.
# Gives Zoom ~2 hours after the 09:05-11:00 IST mandatory session to finish
# processing the AI summary before we pull it.
30 7 * * * samagama /usr/bin/node /var/samagama/server/zoom-fetch-transcripts.js >> /var/samagama/server/logs/zoom-transcripts.log 2>&1
# Catchup verbatim transcript fetch at 13:00 IST (07:30 UTC) daily.
# Same 2-hour buffer for Zoom to finish processing the AI Companion VTT transcript.
# Uses --days 2 to also catch any session that missed the 11:15 IST pipeline run.
31 7 * * * samagama /usr/bin/node /var/samagama/server/zoom-ingest-all-transcripts.js --days 2 >> /var/samagama/server/logs/zoom-transcripts.log 2>&1
